# c100p.service
[Unit]
Description=C100P tweaks

[Service]
User=%I
ExecStart=/usr/bin/c100p.tablet-mode

[Install]
WantedBy=multi-user.target